# rpi-speaker-zones.service
[Install]
WantedBy=multi-user.target

[Service]
ExecStart=/usr/local/bin/rpi-speaker-zones
